# github-auto-update/deploy/github-auto-update.timer.example
[Unit]
Description=Run GitHub multi-project updater every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=github-auto-update.service

[Install]
WantedBy=timers.target
