# pi-status/deploy/pi-health-watchdog.timer.example
[Unit]
Description=Run Raspberry Pi health watchdog every 5 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
Unit=pi-health-watchdog.service

[Install]
WantedBy=timers.target
